# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
logfile_maxbytes=10MB
loglevel=info
user=root

; ── 1. Virtual framebuffer ────────────────────────────────────
; Resolution patched at startup from CORRIDORKEY_RESOLUTION in .env
[program:xvfb]
command=Xvfb :1 -screen 0 CORRIDORKEY_RESOLUTION_PLACEHOLDER -ac +extension GLX +render -noreset
autorestart=true
priority=10
stdout_logfile=/var/log/xvfb.log
stderr_logfile=/var/log/xvfb.log

; ── 2. x11vnc ─────────────────────────────────────────────────
[program:x11vnc]
command=x11vnc -display :1 -forever -shared -nopw -xkb -ncache 10 -rfbport 5900
autorestart=true
priority=20
startsecs=2
startretries=5
stdout_logfile=/var/log/x11vnc.log
stderr_logfile=/var/log/x11vnc.log

; ── 3. noVNC ──────────────────────────────────────────────────
[program:novnc]
command=websockify --web=/usr/share/novnc 6080 localhost:5900
autorestart=true
priority=30
startsecs=3
stdout_logfile=/var/log/novnc.log
stderr_logfile=/var/log/novnc.log

; ── 4. File browser ───────────────────────────────────────────
[program:filebrowser]
command=/usr/local/bin/filebrowser --port 6081 --address 0.0.0.0 --root /app/ClipsForInference --noauth
autorestart=true
priority=35
startsecs=2
stdout_logfile=/var/log/filebrowser.log
stderr_logfile=/var/log/filebrowser.log

; ── 5. EZ-CorridorKey GUI ─────────────────────────────────────
[program:corridorkey]
command=/app/.venv/bin/python /app/main.py --gui
directory=/app
environment=DISPLAY=":1",QT_QPA_PLATFORM="xcb",LIBGL_ALWAYS_SOFTWARE="0",HF_HUB_DISABLE_PROGRESS_BARS="1"
autorestart=true
priority=40
startsecs=5
startretries=3
stdout_logfile=/var/log/corridorkey.log
stderr_logfile=/var/log/corridorkey.log
